# debian/calaos-server.service
[Unit]
Description=Calaos-Server Daemon
After=network.target

[Service]
ExecStart=/usr/bin/calaos-server

[Install]
WantedBy=multi-user.target